# var/service/ps_manager/run
#!/bin/bash
exec /usr/local/bin/pmbus_pwrsply_mngr -d -s 5
