# var/service/qpi_loggerd/run
! /bin/bash

APP="qpi_logger"
FLASHLOCATION="/opt/flash/cisco/bin/${APP}"
INITRAMFSLOCATION="/opt/cisco/bin/${APP}"

if [ -e ${INITRAMFSLOCATION} ]; then
	exec ${INITRAMFSLOCATION};
elif [ -e ${FLASHLOCATION} ]; then
	exec ${FLASHLOCATION};
fi

