# var/service/DimmBL-DbSyncTask/run
#!/bin/bash
exec /usr/local/bin/DimmBL-DbSyncTask
