# var/service/appweb/run
#!/bin/sh
/etc/init.d/appwebd initservice
exec /usr/local/bin/appweb.sh 2>&1 > /dev/null
