# var/service/forkexecproxyd/run
#!/bin/bash
exec /usr/local/bin/forkexecproxyd 1> /dev/null 2> /dev/null
