# var/service/GDSCommDmn/run
#!/bin/bash
exec  /usr/local/bin/GDSCommDmn 10
