# var/service/gpu_mgr/run
#!/bin/bash
exec /usr/local/bin/gpu_manager -d -s 5
