# var/service/llbl/run
#!/bin/bash
exec /usr/local/bin/llbl -fg

