# var/service/video_triggerd/run
#!/bin/bash
exec /usr/local/bin/video_trigger
