# var/service/pof_monitor/run
#!/bin/bash
exec /usr/local/bin/pof_monitord monitor
