# var/service/ibmc_manager/run
#!/bin/bash
sleep 15
exec /usr/local/bin/bmc_manager service
