# var/service/snmpd/run
#! /bin/bash

exec snmpd -f -C -c /var/snmp/snmpd.conf -LS 0-4 d 
