# var/service/mem_throttled/run
! /bin/bash

APP="mem_throttle"
FLASHLOCATION="/opt/flash/cisco/bin/${APP}"
INITRAMFSLOCATION="/opt/cisco/bin/${APP}"

if [ -e ${INITRAMFSLOCATION} ]; then
	exec ${INITRAMFSLOCATION};
elif [ -e ${FLASHLOCATION} ]; then
	exec ${FLASHLOCATION};
fi

