# var/service/rack_monitor/run
#!/bin/bash
exec /usr/local/bin/rack_monitor

