# var/service/vkcs_mux/run
#!/bin/bash
exec /usr/local/bin/vkcs_mux
