# var/service/selparser/run
#!/bin/bash
exec /usr/local/bin/selparser -d> /dev/null 2> /dev/null
