# var/service/ciscovmd/run
#!/bin/bash
exec /usr/local/bin/ciscovm
