# var/service/sam_master/run
#!/bin/bash
exec /usr/local/bin/sam_master

