# var/service/ntpd/run
#!/bin/bash
exec ntpd -g -n 2>&1 /var/log/ntpd.log 

